# ee.conf
[program:ee]
command=python3 /home/curator2/ee/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/ee/test.err.log
stdout_logfile=/home/curator2/ee/test.out.log